# time/test.cron
*/2 * * * * date > /home/pi/time/log.txt
